# systemd/avellcc-restore.service
[Unit]
Description=Restore avellcc keyboard and lightbar state
After=multi-user.target systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/avellcc keyboard --restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
