# systemd/avellcc-restore.service
[Unit]
Description=Restore Avell Control Center keyboard and lightbar state
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/avellcc keyboard --restore
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
